# ops/voice-interview-backend.service
[Unit]
Description=Voice Interview Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/voice-interview/backend
EnvironmentFile=/opt/voice-interview/backend/voice-interview-backend.env
ExecStart=/usr/bin/java -jar /opt/voice-interview/backend/voice-interview-backend-0.0.1-SNAPSHOT.jar --server.port=8080
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
